# src/build/modules/memcached/files/etc/supervisor/conf.d/memcached.conf
[program:memcached]
command = /usr/local/bin/memcached -u root -m 512
autostart = true
autorestart = true
