# src/build/modules/memcached/files/etc/supervisor/conf.d/memcached.conf
[program:memcached]
command=memcached -u root
autostart=true
autorestart=true
