# src/memcached/run/modules/memcached/templates/memcached.conf.erb
[program:memcached]
command = /usr/local/bin/memcached -u root -m <%= @cachesize %> -c <%= @max_connections %> -t <%= @threads %> -I <%= @max_item_size %> -<%= @verbosity %>
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
