# src/run/modules/memcached/templates/memcached.conf.erb
[program:memcached]
command = /usr/local/bin/memcached -u root -m <%= @cachesize %>
autostart = true
autorestart = true
